# system/etc/systemd/system/fstrim.timer.d/override.conf
[Unit]
Description=Run fstrim on every boot and daily at 04:00 localtime.

[Timer]
OnBootSec=0
OnCalendar=*-*-* 04:00:00
Persistent=True
RandomizedDelaySec=0

[Install]
WantedBy=timers.target
