# system/etc/systemd/system/paru-clean.service
[Unit]
Description=Clean paru AUR cache

[Service]
Type=oneshot
ExecStart=/usr/bin/paru -Sc --noconfirm

[Install]
WantedBy=multi-user.target
