# etc/systemd/system/paccache.timer.d/override.conf
[Unit]
Description=Clear pacman cache on every boot and everytime at 04:00 localtime.

[Timer]
OnBootSec=0
OnCalendar=*-*-* 04:00:00
Persistent=true
AccuracySec=1us

[Install]
WantedBy=timers.target
