# home/.config/systemd/user/mic-volume.service
[Unit]
Description=Force mic volume to 0.20 after PipeWire/WirePlumber
After=pipewire.service pipewire-pulse.service wireplumber.service
PartOf=pipewire.service pipewire-pulse.service wireplumber.service

[Service]
Type=oneshot
ExecStart=%h/.local/bin/set-mic-volume.sh

[Install]
WantedBy=default.target
