# home/.config/systemd/user/yay-aur-clean.service
[Unit]
Description=Clean yay AUR cache

[Service]
Type=oneshot
ExecStart=/usr/bin/yay -Sc --noconfirm --aur
